# .config/systemd/user/sommelier-x@0.service.d/cros-sommelier-x-override.conf
[Service]
Environment="SOMMELIER_FRAME_COLOR=#F2F2F2"
Environment="SOMMELIER_ACCELERATORS=Super_L,<Alt>bracketright,<Alt>bracketleft,<Alt>minus,<Alt>equal"
